# provisioning/systemd/chat-embed-snapshots.timer
[Unit]
Description=Daily run of chat-embed-snapshots.service (02:30 Europe/Paris)

[Timer]
OnCalendar=*-*-* 02:30:00 Europe/Paris
Persistent=true
Unit=chat-embed-snapshots.service

[Install]
WantedBy=timers.target
