# provisioning/systemd/chat-embed-snapshots.service.in
[Unit]
Description=Embed yesterday's Bluesky top-10 snapshot into pgvector
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=@PROJECT_DIR@/provisioning/containers
EnvironmentFile=@PROJECT_DIR@/.env.local
ExecStart=/bin/bash -c 'TARGET_DATE=$(TZ=Europe/Paris date -d "yesterday" +%%F); /usr/bin/docker compose -f docker-compose.yaml -f docker-compose.override.yaml exec -T app bin/console chat:embed-snapshots --date="$TARGET_DATE" --no-interaction'
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
